# Yocto/poky/meta-mylayer/recipes-example/service/files/vsomeip.service
[Unit]
Description = vsomeip startup script

[Service]
ExecStart = /usr/bin/example/vsomeip-example/CANService

[Install]
WantedBy = multi-user.target
